# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.


## Set TRX Config during startup
#/opt/fmpoland/sa818/818cli-prog

#no changes required - all deprecated in 1.19.02.2
#restore sound card settings
#alsactl restore


exit 0
