# etc/cron.d/svxlink
CRONDARGS=-1 -m off
MAILTO=""
# Restart svxlink every day at 6:00 and stop on 23:30
30 23 * * * root systemctl stop svxlink
0 6 * * * root systemctl start svxlink
